# files/niri-settings/etc/profile.d/zodium-icons.fish
if set -q __ETC_PROFILE_ZODIUM_SOURCED; return; end
set -x __ETC_PROFILE_ZODIUM_SOURCED 1

set -x XDG_DATA_DIRS $XDG_DATA_DIRS:/etc/zodium