# files/niri-settings/usr/lib/systemd/system/papirus-accent.path
[Unit]
Description=Watch Noctalia theming for accent changes

[Path]
PathExistsGlob=/home/*/.cache/noctalia/theming.dynamic.toml
PathExistsGlob=/home/*/.cache/noctalia/theming.predefined.toml
PathExistsGlob=/home/*/.cache/noctalia/predefined-scheme.json
Unit=papirus-accent.service

[Install]
WantedBy=multi-user.target