# files/niri-settings/usr/lib/systemd/system/papirus-accent.service
[Unit]
Description=Sync Papirus folder color to Noctalia accent

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 5
ExecStart=/usr/libexec/papirus-sync-systemd