# files/niri-settings/usr/lib/systemd/system/papirus-accent.service
[Unit]
Description=Sync Papirus folder color to Noctalia accent
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/libexec/papirus-sync-systemd
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=graphical.target