# files/niri-settings/usr/lib/systemd/system/papirus-accent.path
[Unit]
Description=Watch Noctalia theming for accent changes

[Path]
PathChanged=%h/.cache/noctalia/theming.dynamic.toml
Unit=papirus-accent.service

[Install]
WantedBy=default.target