# files/niri-settings/usr/lib/systemd/system/papirus-accent.service
[Unit]
Description=Sync Papirus folder color to Noctalia accent
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/libexec/papirus-sync-systemd
Restart=on-failure

[Install]
WantedBy=multi-user.target